# files/systemd/backup_private_files.service
[Unit]
Description=Backup important configuration files

[Service]
Type=oneshot
ExecStart=/home/freeman/.scripts/backup_private_files.sh
